# saemref/files/saemref-supervisor.conf
{% from "saemref/map.jinja" import saemref with context -%}
[program:saemref]
{% if saemref.instance.wsgi -%}
command=uwsgi --ini /home/{{ saemref.instance.user }}/etc/cubicweb.d/{{ saemref.instance.name }}/uwsgi.ini
{% else -%}
command=/home/{{ saemref.instance.user }}/venv/bin/cubicweb-ctl start -D {{ saemref.instance.name }}
{% endif -%}
user={{ saemref.instance.user }}
redirect_stderr=true
stdout_logfile=/home/{{ saemref.instance.user }}/{{ saemref.instance.name }}.log
environment=CW_MODE=user,HOME=/home/{{ saemref.instance.user }}
numprocs=1
autostart=true
autorestart=true
startsecs=10
