# saemref/files/saemref-supervisor.conf
{% from "saemref/map.jinja" import saemref with context -%}
[program:saemref]
{% if saemref.instance.wsgi -%}
command=/home/{{ saemref.instance.user }}/venv/bin/gunicorn 'cubicweb.pyramid:wsgi_application()' --forwarded-allow-ips * --threads {{ saemref.instance.wsgi_threads }} --workers {{ saemref.instance.wsgi_workers }} --bind 0.0.0.0:{{ saemref.instance.port }}
{% else -%}
command=/home/{{ saemref.instance.user }}/venv/bin/cubicweb-ctl pyramid -D {{ saemref.instance.name }}
{% endif -%}
user={{ saemref.instance.user }}
redirect_stderr=true
stdout_logfile=/home/{{ saemref.instance.user }}/{{ saemref.instance.name }}.log
environment=CW_MODE=user,HOME=/home/{{ saemref.instance.user }},CW_INSTANCE={{ saemref.instance.name }}
autorestart=true
startsecs=10

[program:saemref-oai]
{% if saemref.instance.wsgi -%}
command=/home/{{ saemref.instance.user }}/venv/bin/gunicorn 'cubicweb.pyramid:wsgi_application()' --forwarded-allow-ips * --threads {{ saemref.instance.oai_threads }} --workers 1 --bind 0.0.0.0:{{ saemref.instance.oai_port }}
{% else -%}
command=/home/{{ saemref.instance.user }}/venv/bin/cubicweb-ctl pyramid -D {{ saemref.instance.name }}
{% endif -%}
environment=CW_MODE=user,HOME=/home/{{ saemref.instance.user }},CW_INSTANCE={{ saemref.instance.name }},CW_CONNECTIONS_POOL_SIZE={{ saemref.instance.oai_threads }},CW_PORT={{ saemref.instance.oai_port }}
user={{ saemref.instance.user }}
redirect_stderr=true
stdout_logfile=/home/{{ saemref.instance.user }}/{{ saemref.instance.name }}-oai.log
autorestart=true
startsecs=10
