# saemref/files/saemref-supervisor.conf
{% from "saemref/map.jinja" import saemref with context -%}
[program:saemref]
{% if saemref.instance.wsgi -%}
command=/home/{{ saemref.instance.user }}/venv/bin/gunicorn 'cubicweb.pyramid:wsgi_application()' --threads {{ saemref.instance.wsgi_threads }} --workers {{ saemref.instance.wsgi_workers }} --bind 0.0.0.0:{{ saemref.instance.port|default('8080', true) }}
{% else -%}
command=/home/{{ saemref.instance.user }}/venv/bin/cubicweb-ctl pyramid -D {{ saemref.instance.name }}
{% endif -%}
user={{ saemref.instance.user }}
redirect_stderr=true
stdout_logfile=/home/{{ saemref.instance.user }}/{{ saemref.instance.name }}.log
environment=CW_MODE=user,HOME=/home/{{ saemref.instance.user }},CW_INSTANCE={{ saemref.instance.name }}
numprocs=1
autostart=true
autorestart=true
startsecs=10
